# rootfs/etc/systemd/system/glyphguard-watchdog.service
[Unit]
Description=🧠 GlyphGuard Auto-Lock Watchdog Daemon
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/trinity/glyphguard/glyphguard_watchdog.py
Restart=on-failure
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target