# rootfs/etc/systemd/system/trinity_boot_injector.service
[Unit]
Description=Trinity Boot Injector
After=multi-user.target

[Service]
Type=oneshot
ExecStart=/opt/trinity/bin/trinity_boot_injector.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
